# cron.d/bot-cron
*/1 * * * * /usr/local/bin/python /usr/src/app/bot.py >> /var/log/cron.log 2>&1
